# server/planos-api.container
# Quadlet — copiar a ~/.config/containers/systemd/planos-api.container (usuario 'planos')
# Recargar:  systemctl --user daemon-reload && systemctl --user start planos-api
[Unit]
Description=SAF Planos API (Express + node-oracledb)
After=network-online.target

[Container]
Image=localhost/planos-api:1.0
# Solo se publica en loopback del host; Nginx (host) hace de fachada TLS
PublishPort=127.0.0.1:3000:3000

# Credenciales Oracle como secrets (NO en texto plano)
Secret=oracle_user,type=mount
Secret=oracle_password,type=mount
Environment=ORACLE_CONNSTR=oracle-db.aicsa.local/SAFDB
Environment=XFDF_USE_DB=true
Environment=TRUST_PROXY=true

# Persistencia de XFDF en disco del host (fallback al CLOB de Oracle)
Volume=/srv/planos/xfdf:/data/xfdf:Z

# ── Endurecimiento ──────────────────────────────────────────────────────
ReadOnly=true
Tmpfs=/tmp:rw,noexec,nosuid
NoNewPrivileges=true
DropCapability=ALL
PidsLimit=200
Memory=512m

[Service]
Restart=always
TimeoutStartSec=60

[Install]
WantedBy=default.target
